# deploy/supervisord.conf
; EnlargeImage backend supervisord config
; Drop into /etc/supervisor/conf.d/enlargeimage.conf on the target host
; (or include via [include] in supervisord.conf).
;
; Usage:
;   sudo supervisorctl reread
;   sudo supervisorctl update
;   sudo supervisorctl start enlargeimage
;   sudo supervisorctl status enlargeimage

[program:enlargeimage]
; Run from the backend dir so .env, models/, storage/ paths resolve.
directory=/opt/enlargeimage/backend
command=/opt/enlargeimage/backend/.venv/bin/python run.py
                    --host 0.0.0.0
                    --port 8000
                    --model-dir /opt/enlargeimage/backend/models
                    --storage-dir /opt/enlargeimage/backend/storage
                    --db-path /opt/enlargeimage/backend/data.db
                    --workers 1
                    --log-level info

; Run as a dedicated service user (create first: useradd -r -s /bin/false enlargeimage)
user=enlargeimage

autostart=true
autorestart=true
startsecs=10
startretries=3
stopwaitsecs=30
stopsignal=TERM

; Redirect logs to dedicated files; supervisord rotates them via logrotate.
stdout_logfile=/var/log/enlargeimage/enlargeimage.out.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=10
stderr_logfile=/var/log/enlargeimage/enlargeimage.err.log
stderr_logfile_maxbytes=50MB
stderr_logfile_backups=10

; Production hardening
environment=PATH="/opt/enlargeimage/backend/.venv/bin:/usr/bin:/bin"
environment=PYTHONUNBUFFERED="1"
environment=PYTHONDONTWRITEBYTECODE="1"
priority=100
